# trabajo.cron
*/10 * * * * /home/USUARIO/usuarios_ldap/usuarios.py
